# vermeer/config/supervisor.conf
<!--
  Licensed to the Apache Software Foundation (ASF) under one or more
  contributor license agreements.  See the NOTICE file distributed with
  this work for additional information regarding copyright ownership.
  The ASF licenses this file to You under the Apache License, Version 2.0
  (the "License"); you may not use this file except in compliance with
  the License.  You may obtain a copy of the License at

      http://www.apache.org/licenses/LICENSE-2.0

  Unless required by applicable law or agreed to in writing, software
  distributed under the License is distributed on an "AS IS" BASIS,
  WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
  See the License for the specific language governing permissions and
  limitations under the License.
  -->
  
[inet_http_server]
port=:8999

[supervisord]
logfile=./supervisord.log
logfile_maxbytes=10485760
logfile_backups=3

[program-default]
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=3
stderr_logfile_maxbytes=10485760
stderr_logfile_backups=3
restart_when_binary_changed=true
restartpause=10
autorestart=true
autostart=true

[program:vermeer_master]
command = ./vermeer.sh launch master
process_name=vermeer_master
stdout_logfile=./master.log
stderr_logfile=./master_err.log

[program:vermeer_worker]
command = ./vermeer.sh launch worker
process_name=vermeer_worker
stdout_logfile=./worker.log
stderr_logfile=./worker_err.log

[program:mem_supervisor]
command = sh mem_supervisor.sh
process_name = mem_supervisor
stdout_logfile=./mem_supervisor.log
stderr_logfile=./mem_supervisor_err.log
